# .docker/cron.daily/welcome-update.sh
#!/bin/sh

cd /var/www/osm-welcome-tool/
php bin/console welcome:update
